# rc.local
--- insert before 'exit 0' --- remove this line ---

# Part of S0/Impulse to Volkszaehler 'RaspberryPI deamon'
echo 17 > /sys/class/gpio/export && echo "in" > /sys/class/gpio/gpio17/direction && echo "rising" > /sys/class/gpio/gpio17/edge
echo 18 > /sys/class/gpio/export && echo "in" > /sys/class/gpio/gpio18/direction && echo "rising" > /sys/class/gpio/gpio18/edge
echo 22 > /sys/class/gpio/export && echo "in" > /sys/class/gpio/gpio22/direction && echo "rising" > /sys/class/gpio/gpio22/edge
echo 23 > /sys/class/gpio/export && echo "in" > /sys/class/gpio/gpio23/direction && echo "rising" > /sys/class/gpio/gpio23/edge
echo 24 > /sys/class/gpio/export && echo "in" > /sys/class/gpio/gpio24/direction && echo "rising" > /sys/class/gpio/gpio24/edge
echo 27 > /sys/class/gpio/export && echo "in" > /sys/class/gpio/gpio27/direction && echo "rising" > /sys/class/gpio/gpio27/edge
